# docker/supervisor.worker.ini
[supervisord]
nodaemon=true

[program:celery]
user=pyazo
command=celery worker --autoscale=10,3 -E -B -A=pyazo.root.celery -s=/tmp/celerybeat-schedule
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
